# package/video_webapp/video_webapp.service
[Unit]
Description=Video Webapp

[Service]
ExecStart=/usr/bin/videowebapp
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target